# LaunchAgents/com.technicalpickles.qmd-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.technicalpickles.qmd-refresh</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>$HOME/.bun/bin/qmd update &amp;&amp; $HOME/.bun/bin/qmd embed</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/com.technicalpickles.qmd-refresh.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/com.technicalpickles.qmd-refresh.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
